# scripts/maxos.service
[Unit]
Description=MaxOS API Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/maxos
EnvironmentFile=-/etc/maxos.env
ExecStart=/usr/bin/env uvicorn max_os.interfaces.api.main:app --host 0.0.0.0 --port 8000
Restart=on-failure

[Install]
WantedBy=multi-user.target
